# xci/source/info.apps/Manage-XDCDB-Resources/tags/tag-1.4-20170331/sbin/speedpage.cron
#!/bin/sh
LOCATION=/soft/warehouse-apps-1.0/Manage-Speedpage/
PRODLOC=$LOCATION/PROD/
CACHE=$LOCATION/var

DATE=`date +%Y-%m-%dT%H:%M:%SZ`
#LOG=/soft/info/logs/outages.log
LOG=$CACHE/speedpage.log
echo $DATE >> $LOG
$PRODLOC/sbin/cache.speed.pl --cache=$CACHE >> $LOG 2>&1

$PRODLOC/sbin/route_speedpage.sh
