# xci/source/info.apps/Manage-XDCDB-Resources/tags/tag-1.9-20201026b/sbin/V4XDCDB.cron
#!/bin/sh
APP_BASE=/soft/warehouse-apps-1.0/Manage-XDCDB-Resources
APP_SOURCE=${APP_BASE}/PROD
APP_CACHE=${APP_BASE}/var
${APP_SOURCE}/sbin/cache.xdcdb.pl --cache=${APP_CACHE}

DATE=`date +%Y-%m-%dT%H:%M:%SZ`
LOG=${APP_CACHE}/outages.log

${APP_SOURCE}/sbin/route_xdcdb.sh ${APP_CACHE}/tgresources.csv
